# HW4/task4.cron
#Задание 4: Создать файл crontab, который ежедневно регистрирует занятое каждым пользователем дисковое пространство в его домашней директории. 
SHELL=/bin/sh

0 0 * * * /home/maxim/homework/linux_basics/HW4/task4_shell.sh 


